# shared/systemd_user/fox-monitor-watch.service
[Unit]
Description=React to Hyprland monitor hot-swap — rerun rotate / personalise / re-pin
After=default.target
# graphical-session.target is RefuseManualStart=yes upstream and Hyprland
# does not pull it in. WantedBy=default.target activates the unit on user
# login, which is what we want for a long-running Hyprland-aware script.

[Service]
Type=simple
ExecStart=%h/.config/hypr/scripts/fox-monitor-watch.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
