# shared/systemd_user/foxml-focus-pulse.service
[Unit]
Description=Brief OSD on Hyprland workspace switch (foxml focus-pulse)
After=default.target
# graphical-session.target is RefuseManualStart=yes upstream and Hyprland
# does not pull it in; default.target activates on user login.

[Service]
Type=simple
ExecStart=%h/.config/hypr/scripts/focus-pulse.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
