# shared/systemd_user/power-state-watcher.service
[Unit]
Description=Swap hypridle config based on AC/battery power state
After=default.target
# graphical-session.target is RefuseManualStart=yes upstream and Hyprland
# does not pull it in; default.target activates on user login.

[Service]
Type=simple
ExecStart=%h/.config/hypr/scripts/power_state_watcher.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
